# buildroot/board/cvitek/SG200X/overlay/etc/profile.d/middleware.sh
export PATH=/mnt/system/usr/bin:$PATH
export LD_LIBRARY_PATH=/root/libs_patch/lib:/root/libs_patch/middleware_v2:/root/libs_patch/middleware_v2_3rd:/root/libs_patch/tpu_sdk_libs:/root/libs_patch:/root/libs_patch/opencv:$LD_LIBRARY_PATH
